# debian/maas.postrm
#!/bin/sh -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ] ; then

	# Handle custom squid removal
	if [ -e /etc/squid3/squid.conf.orig ] ; then
		if [ -h /etc/squid3/squid.conf ]; then
			rm -f /etc/squid3/squid.conf
		fi
		cp /etc/squid3/squid.conf.orig /etc/squid3/squid.conf
	fi

	invoke-rc.d squid3 restart

fi

#DEBHELPER#
exit 0
