# debian/maas.upstart
# maas - provisioning service
#
# MaaS Provisioning Service

description     "MaaS"
author "Andres Rodriguez <andres.rodriguez@canonical.com>"

start on filesystem and net-device-up
stop on runlevel [016]

respawn

# To add options to your daemon, edit the line below:
exec /usr/bin/twistd -n maas-pserv --config-file=/etc/maas/pserv.yaml
