# debian/maas.maas-pserv.upstart
# maas-pserv - provisioning service
#
# MAAS Provisioning Service

description     "MAAS"
author "Andres Rodriguez <andres.rodriguez@canonical.com>"

start on filesystem and net-device-up
stop on runlevel [016]

respawn

# To add options to your daemon, edit the line below:
exec /usr/bin/twistd -n --uid=maas --gid=maas --pidfile=/run/maas-pserv.pid --logfile=/dev/null maas-pserv --config-file=/etc/maas/pserv.yaml
