# debian/maas-region-controller.prerm
#!/bin/sh

set -e

. /usr/share/debconf/confmodule

#DEBHELPER#

if [ -x /usr/sbin/invoke-rc.d ]; then
	invoke-rc.d apache2 stop || true
else
	/etc/init.d/apache2 stop || true
fi

. /usr/share/dbconfig-common/dpkg/prerm.pgsql
dbc_go maas-region-controller $@

exit 0
