# debian/maas-proxy.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ]
then
    # Stop the squid3 instance that came with the isc-dhcp-server package.
    # We run our own squid3 instances.
    invoke-rc.d squid3 stop

    # Ensure log folder is created.
    mkdir -p /var/log/maas/proxy
    chown -R proxy:proxy /var/log/maas/proxy

fi

#DEBHELPER#
