# debian/maas-region-controller-min.maas-regiond.upstart
description "MAAS Region Controller"
author "Gavin Panella <gavin.panella@canonical.com>"

# This starts multiple instances of `maas-regiond-worker` and does
# nothing else itself; there is no main script so this runs "forever".
# By default, it launches four workers.

start on filesystem and net-device-up
stop on runlevel [016]

respawn

pre-start script
    for core in $(seq 4); do
        start maas-regiond-worker core=${core}
    done
end script
